# etc/mb-mail-service.socket
[Unit]
Description=Mail Service Socket

[Socket]
ListenStream=0.0.0.0:3000

[Install]
WantedBy=sockets.target
